# etc/supervisor.conf
[group:crypto]
programs=crypto-trader,crypto-supervisor

[program:crypto-trader]
directory=/home/%USERNAME%/current
command=/home/%USERNAME%/venv/bin/python app/trader.py
user=%USERNAME%
autorestart=true
autostart=true
redirect_stderr=true
stdout_logfile=/home/%USERNAME%/logs/supervisor-crypto-profit-orders.log
stdout_logfile_maxbytes=10MB
stdout_capture_maxbytes=10MB
stdout_logfile_backups=10

[program:crypto-supervisor]
directory=/home/%USERNAME%/current
command=/home/%USERNAME%/venv/bin/python app/supervisor.py
user=%USERNAME%
autorestart=true
autostart=true
redirect_stderr=true
stdout_logfile=/home/%USERNAME%/logs/supervisor-crypto-init-orders.log
stdout_logfile_maxbytes=10MB
stdout_capture_maxbytes=10MB
stdout_logfile_backups=10