# etc/supervisor.conf
[group:crypto]
programs=crypto-trader

[program:crypto-trader]
directory=/home/%USERNAME%/current
command=/home/%USERNAME%/venv/bin/python app/trader.py
user=%USERNAME%
stopsignal=INT
autorestart=true
autostart=true
redirect_stderr=true
stdout_logfile=/home/%USERNAME%/logs/supervisor.txt
stdout_logfile_maxbytes=1MB
stdout_capture_maxbytes=1MB
stdout_logfile_backups=100
